# rundeck-launcher/rpm/src/main/resources/scripts/postinst
if [ ! -e ~${rpmAppUser}/.ssh/id_rsa ]; then
	su -c "ssh-keygen -q -t rsa -C '' -N '' -f ~${rpmAppUser}/.ssh/id_rsa" ${rpmAppUser}
fi
/sbin/chkconfig --add rundeckd
